# libs/supervisor/supervisord.d/ai.conf
[program:ai]
directory=/ai/PDF-OCR-RTP
command=bash /etc/supervisord.d/celery.sh
user=karen
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 5
stdout_logfile = /var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile = /var/log/supervisor/%(program_name)s-stderr.log
